# deploy/generation-workers.supervisor.conf
[program:ai_creator_submit_worker]
command=/bin/bash /www/wwwroot/ai-creator/deploy/run_generation_worker.sh submit
directory=/www/wwwroot/ai-creator
autostart=true
autorestart=true
startsecs=3
startretries=5
stderr_logfile=/www/wwwlogs/ai_creator_submit_worker.err.log
stdout_logfile=/www/wwwlogs/ai_creator_submit_worker.out.log
user=root
stopasgroup=true
killasgroup=true

[program:ai_creator_poll_worker]
command=/bin/bash /www/wwwroot/ai-creator/deploy/run_generation_worker.sh poll
directory=/www/wwwroot/ai-creator
autostart=true
autorestart=true
startsecs=3
startretries=5
numprocs=3
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/www/wwwlogs/ai_creator_poll_worker_%(process_num)02d.err.log
stdout_logfile=/www/wwwlogs/ai_creator_poll_worker_%(process_num)02d.out.log
user=root
stopasgroup=true
killasgroup=true

[program:ai_creator_oss_worker]
command=/bin/bash /www/wwwroot/ai-creator/deploy/run_generation_worker.sh sync
directory=/www/wwwroot/ai-creator
autostart=true
autorestart=true
startsecs=3
startretries=5
stderr_logfile=/www/wwwlogs/ai_creator_oss_worker.err.log
stdout_logfile=/www/wwwlogs/ai_creator_oss_worker.out.log
user=root
stopasgroup=true
killasgroup=true
